# supervisord.conf
[program:xrdp]
command=/usr/sbin/xrdp -nodaemon
process_name = xrdp

[program:xrdp-sesman]
command=/usr/sbin/xrdp-sesman --nodaemon
process_name = xrdp-sesman
